# client/simburd.conf
author "Larry Reid <lcreid@jadesystems.ca>
description "Start the SIMple BackUp and Recovery daemon"

start on started networking

exec /usr/bin/simburd

