# client/debian/DEBIAN/postinst
#!/bin/sh
echo Running postinst
systemctl enable simburd
